# oa_sourcing_hub_clean/rootfs/etc/services.d/oa/run
#!/bin/sh
set -eu
echo "[oa][run] service starting (clean 1.5.2)"
cd /app
exec /opt/venv/bin/streamlit run app.py --server.port 8501 --server.address 0.0.0.0
