# oa_sourcing_hub_ultra/rootfs/etc/services.d/oahub/run
#!/bin/sh
set -eu
echo "[oahub] service starting (pure /bin/sh)"
cd /app
exec /opt/venv/bin/streamlit run app.py --server.port 8501 --server.address 0.0.0.0
