# s6_pure_test/rootfs/etc/services.d/s6pure/run
#!/bin/sh
set -eu
echo "[s6pure] s6 service started ✅ (pure /bin/sh)"
exec tail -f /dev/null
