# s6_zero_test/rootfs/etc/services.d/oauc/run
#!/bin/sh
set -eu
echo "[s6-zero] s6 service started ✅ (pure /bin/sh)"
exec tail -f /dev/null
