# hello_s6_test/rootfs/etc/services.d/hello/run
#!/bin/sh
set -eu
echo "[hello] s6 service started ✅"
# keep container alive
exec tail -f /dev/null
